# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx

[program:holo]
command=npm run dev --workspace holo
directory=/usr/src/app